# deployment/transcriptionbot-django.service
[Unit]
Description=TranscriptionBot Django Web Application (Gunicorn)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=notify
User=company
Group=company
WorkingDirectory=/opt/applications/projects/TranscriptionBot/TranscriptionBot/django_admin
Environment="PATH=/opt/applications/projects/TranscriptionBot/TranscriptionBot/.venv/bin"
Environment="PYTHONPATH=/opt/applications/projects/TranscriptionBot/TranscriptionBot/"
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=/opt/applications/projects/TranscriptionBot/TranscriptionBot/.env

ExecStart=/opt/applications/projects/TranscriptionBot/TranscriptionBot/.venv/bin/gunicorn \
    --config /opt/applications/projects/TranscriptionBot/TranscriptionBot/deployment/gunicorn_config.py \
    --workers 3 \
    --bind 127.0.0.1:8100 \
    config.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transcriptionbot-django

# Security and resource limits
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/applications/projects/TranscriptionBot/TranscriptionBot/django_admin/media /opt/applications/projects/TranscriptionBot/TranscriptionBot/django_admin/staticfiles /opt/applications/projects/TranscriptionBot/TranscriptionBot/django_admin/logs /opt/applications/projects/TranscriptionBot/TranscriptionBot/logs /var/run/transcriptionbot
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
