# deployment/transcriptionbot-telegram.service
[Unit]
Description=TranscriptionBot Telegram Bot Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=company
Group=company
WorkingDirectory=/opt/applications/projects/TranscriptionBot/TranscriptionBot/bot
Environment="PATH=/opt/applications/projects/TranscriptionBot/TranscriptionBot/.venv/bin"
Environment="PYTHONPATH=/opt/applications/projects/TranscriptionBot/TranscriptionBot/"
EnvironmentFile=/opt/applications/projects/TranscriptionBot/TranscriptionBot/.env

ExecStart=/opt/applications/projects/TranscriptionBot/TranscriptionBot/.venv/bin/python /opt/applications/projects/TranscriptionBot/TranscriptionBot/bot/main.py

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transcriptionbot-telegram

# Security and resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/applications/projects/TranscriptionBot/TranscriptionBot/media /opt/applications/projects/TranscriptionBot/TranscriptionBot/logs
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
